# systemd/git-watch-blog.service
[Unit]
Description=Git Auto-Commit for Hugo Blog
After=network.target

[Service]
Type=simple
ExecStart=/root/dockeryml/blog/scripts/watch-blog.sh
Restart=always
RestartSec=10
User=root

[Install]
WantedBy=multi-user.target